# backend/systemd/twix-scraper.timer
[Unit]
Description=Twix Saver スクレイピングタイマー
Requires=twix-scraper.service

[Timer]
# 15分間隔で実行（設定により調整可能）
OnCalendar=*:0/15

# 予測可能性を避けるためのランダム遅延（最大2分）
RandomizedDelaySec=120

# システム起動時に実行時刻を逃した場合、起動後に実行
Persistent=true

# タイマーの精度（省電力のため1分に設定）
AccuracySec=1min

[Install]
WantedBy=timers.target