# backend/systemd/twix-scraper.service
[Unit]
Description=Twix Saver
After=network.target
Wants=mongodb.service

[Service]
Type=oneshot
User=your-user
Group=your-group
WorkingDirectory=/path/to/project/backend

# 環境変数の設定
Environment="PATH=/path/to/project/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/path/to/project/backend/.env

# 実行コマンド（データベース駆動でジョブを実行）
ExecStart=/path/to/project/venv/bin/python /path/to/project/backend/main.py --run-jobs

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=twix-scraper

# リソース制限
MemoryMax=2G
CPUQuota=50%

# 障害時の対応
Restart=on-failure
RestartSec=30
TimeoutStartSec=300
TimeoutStopSec=30

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/path/to/project/backend/data /path/to/project/backend/logs

[Install]
WantedBy=default.target